# deploy/enably-admin.service.example
## /etc/systemd/system/enably-admin.service
[Unit]
Description=Enably Admin (Next.js standalone)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/srv/enably-admin/current
EnvironmentFile=/etc/enably-admin.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=15

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
