# deploy/enably-admin.service.example
## /etc/systemd/system/enably-admin.service
[Unit]
Description=Enably Admin (Next.js standalone)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Amazon Linux 2023 default user is ec2-user. Use 'ubuntu' on Ubuntu AMIs.
User=ec2-user
Group=ec2-user
WorkingDirectory=/srv/enably-admin/current
EnvironmentFile=/etc/enably-admin.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=15

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
